# .config/autostart/ssh-agent.sh.desktop
[Desktop Entry]
Exec=~/.config/custom-autostart/ssh-agent.sh
Icon=dialog-scripts
Name=ssh-agent.sh
Path=
Type=Application
X-KDE-AutostartScript=true
